# scripts/systemd/backup-home-vivi.service
[Unit]
Description=Backup /home/vivi and upload to Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
NoNewPrivileges=yes
PrivateTmp=yes
User=vivi
Group=vivi
Environment=HOME=/home/vivi
Environment=RCLONE_CONFIG=/home/vivi/.config/rclone/rclone.conf
Environment=SOURCE_DIR=/home/vivi
Environment=BACKUP_DIR=/home/vivi/.local/share/home_backups
Environment=RCLONE_TARGET=gdrive:vivi-home-backups
# Optional: upload directly into a specific subfolder within RCLONE_TARGET.
# Leave empty to sync directly into RCLONE_TARGET.
Environment=RCLONE_SYNC_PATH=
ExecStart=/usr/bin/env bash /home/vivi/pixelated/scripts/backup/backup-home-vivi.sh
