# scripts/systemd/pixelated-memory-backup.timer.template
[Unit]
Description=Run Pixelated memory DB backup on schedule

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
